# debian/DEBIAN/postinst
#!/bin/sh -e

update-rc.d bautablink defaults
update-rc.d bautaconf defaults
if [ ! -f /bauta/etc/usbconfpw ]; then
	date | md5sum | cut -c -32 >/bauta/etc/usbconfpw
fi
