# debian/DEBIAN/postrm
#!/bin/sh -e

update-rc.d bautablink remove
update-rc.d bautaconf remove

